# rc.local
cd /home/pi/minimeet && startx ./start.sh 
